# gunicorn.service
[Unit]
Description=gunicorn daemon for Dance Portal
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/opt/dance_portal_starter
EnvironmentFile=/opt/dance_portal_starter/.env

# Gunicorn start command
ExecStart=/opt/dance_portal_starter/venv/bin/gunicorn \
          --workers 3 \
          --umask 007 \
          --bind unix:/opt/dance_portal_starter/danceportal.sock \
          dance_portal.wsgi:application

# ✅ Auto-restart settings
Restart=always
RestartSec=5
TimeoutStopSec=30

# Ensure systemd knows how to kill the service cleanly
KillMode=mixed

[Install]
WantedBy=multi-user.target
